# conf/supervisor.conf
[program:instabase]
command=gunicorn -k flask_sockets.worker app:app
directory=/home/ubuntu/instabase
user=nobody
autostart=true
autorestart=true
redirect_stderr=true
